# deploy/telegram-assistant-reminder.service
[Unit]
Description=Telegram Assistant Reminder Worker
Documentation=file:/opt/telegram-assistant/README.md
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=telegram-assistant
Group=telegram-assistant
WorkingDirectory=/opt/telegram-assistant

ExecStart=/opt/telegram-assistant/.venv/bin/python -m telegram_assistant.reminder_worker
TimeoutStartSec=25s

Environment=APP_ENV=production
Environment=PYTHONPATH=/opt/telegram-assistant/src
Environment=APP_DB_PATH=/var/lib/telegram-assistant/events.db
Environment=APP_TIMEZONE=America/Sao_Paulo

LoadCredential=telegram_token:/etc/telegram-assistant/telegram_token
LoadCredential=authorized_user_id:/etc/telegram-assistant/authorized_user_id

StateDirectory=telegram-assistant
StateDirectoryMode=0700
UMask=0077

NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
PrivateTmp=true
PrivateDevices=true
PrivateMounts=true
ProtectSystem=strict
ProtectHome=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectProc=invisible
ProcSubset=pid
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

InaccessiblePaths=-/run/docker.sock
InaccessiblePaths=-/var/run/docker.sock
InaccessiblePaths=-/var/lib/docker
InaccessiblePaths=-/etc/ssh
InaccessiblePaths=-/etc/traefik

MemoryMax=128M
TasksMax=16
LimitNOFILE=64
CPUQuota=20%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-assistant-reminder
