# deploy/telegram-assistant-bot.service
[Unit]
Description=Telegram Assistant Secure Bot
Documentation=file:/opt/telegram-assistant/README.md
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=telegram-assistant
Group=telegram-assistant
WorkingDirectory=/opt/telegram-assistant

ExecStart=/opt/telegram-assistant/.venv/bin/python -m telegram_assistant.bot
Restart=on-failure
RestartSec=10s
TimeoutStartSec=30s
TimeoutStopSec=30s

Environment=APP_ENV=production
Environment=PYTHONPATH=/opt/telegram-assistant/src
Environment=APP_DB_PATH=/var/lib/telegram-assistant/events.db
Environment=APP_TIMEZONE=America/Sao_Paulo
Environment=APP_LLM_AGENDA_DAYS=30
Environment=APP_LLM_AGENDA_MAX_EVENTS=50
Environment=APP_WORKDAY_START=09:00
Environment=APP_WORKDAY_END=18:00
Environment=APP_SLOT_INTERVAL_MINUTES=15
Environment=APP_EVENT_BUFFER_MINUTES=15
Environment=APP_MINIMUM_LEAD_MINUTES=30
Environment=GROQ_RPM_LIMIT=30
Environment=GROQ_RPD_LIMIT=1000
Environment=GROQ_TPM_LIMIT=8000
Environment=GROQ_TPD_LIMIT=200000

LoadCredential=telegram_token:/etc/telegram-assistant/telegram_token
LoadCredential=authorized_user_id:/etc/telegram-assistant/authorized_user_id
LoadCredential=groq_api_key:/etc/telegram-assistant/groq_api_key

StateDirectory=telegram-assistant
StateDirectoryMode=0700
UMask=0077

NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
PrivateTmp=true
PrivateDevices=true
PrivateMounts=true
ProtectSystem=strict
ProtectHome=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectProc=invisible
ProcSubset=pid
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

InaccessiblePaths=-/run/docker.sock
InaccessiblePaths=-/var/run/docker.sock
InaccessiblePaths=-/var/lib/docker
InaccessiblePaths=-/etc/ssh
InaccessiblePaths=-/etc/traefik

MemoryMax=256M
TasksMax=32
LimitNOFILE=128
CPUQuota=25%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-assistant

[Install]
WantedBy=multi-user.target
